# wisp-backend.service
[Unit]
Description=WISP Manager Backend (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
User=laborin
WorkingDirectory=/home/laborin/wisp-manager/backend
EnvironmentFile=/home/laborin/wisp-manager/.env
ExecStart=/home/laborin/wisp-manager/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
